# cx-blog-webhook.service
[Unit]
Description=CX Blog GitHub Webhook Auto Deploy
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/cx-blog
ExecStart=/usr/bin/node /home/ubuntu/cx-blog/webhook-server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cx-blog-webhook

# 安全加固
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
